# packages/core/upstart-init-dev/src/etc/init/cold-plug.conf
description "Cold plug devices"

start on (startup and started udev)

task

exec udevadm trigger
post-stop exec udevadm settle
