# packages/core/upstart-init-dev/src/etc/init/udev.conf
description "Device nodes manager"

start on stopped devfs

expect fork
respawn

exec /sbin/udevd --daemon
