# packages/core/base/src/etc/rc.local
#!/bin/sh

# set printk logging level
echo 8 > /proc/sys/kernel/printk

